# hermie.service
[Unit]
Description=Hermie Temperature Display
After=network.target

[Service]
Type=simple
User=odio
WorkingDirectory=/home/odio/Hacking/hermie
ExecStart=/home/odio/Hacking/hermie/.venv/bin/python /home/odio/Hacking/hermie/watch.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
